# ansible/getraenkekasse/files/getraenkekasse-online.service
[Unit]
Description=Getränkekasse Online
After=docker.service

[Service]
Type=notify
ExecStart=/etc/kalkspace/getraenkekasse/getraenkekasse-online.sh
